# workspace/crontab/laradock
* * * * * laradock php /var/www/b365-api/artisan schedule:run >> /dev/null 2>&1
